# config/systemd/user/audio-precedence-watch.service
[Unit]
Description=CHG-007 audio device precedence follower (sink/source + Rock EQ glue)
After=pipewire.service wireplumber.service pipewire-pulse.service
Wants=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStart=%h/.local/bin/audio-precedence-watch
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
